# ansible/roles/jenkins/files/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/bin/bash -c "rm -rf /run/apache2/* /tmp/httpd* && apachectl -D FOREGROUND"

[program:jenkins]
command=/bin/bash -c "/usr/local/bin/jenkins.sh &"